# nodes/EL9/docker/cron.d/docker-prune
01 01 * * * root docker system prune -f

